# systemd/nas-backup.timer
[Unit]
Description=Backup NAS quotidien (incrémental Btrfs root+home + boot/efi)

[Timer]
# Tous les jours à 02:00. Persistent : si la machine était éteinte à l'heure
# prévue, le backup se lance au prochain démarrage.
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
