# systemd/nas-backup-full.timer
[Unit]
Description=Backup NAS mensuel COMPLET (--full, borne la chaîne d'incréments)

[Timer]
# Le 1er de chaque mois à 03:00 — donc APRÈS le run quotidien de 02:00 : le
# --full purge alors toute la chaîne du mois précédent (y compris l'inc du
# matin même) et repart sur une base propre.
OnCalendar=*-*-01 03:00:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
