# systemd/user/nas-backup-reminder-full.timer
[Unit]
Description=Rappel hebdomadaire de sauvegarde NAS complète (--full)

[Timer]
# Chaque dimanche à 20:30 — un full hebdomadaire borne la chaîne d'incréments
# à ~7 fichiers : restore plus court, moins de points de défaillance.
# (L'ancienne chaîne n'est purgée du NAS qu'après vérification sha256 du
# nouveau full — voir nas-backup.sh.)
OnCalendar=Sun *-*-* 20:30:00
Persistent=true
RandomizedDelaySec=20min

[Install]
WantedBy=timers.target
